# dist/systemd/qiftop-agent.service
[Unit]
Description=qiftop privileged network statistics agent
Documentation=https://github.com/ines/qiftop

[Service]
Type=dbus
BusName=org.qiftop.NetworkAgent1
ExecStart=/usr/bin/qiftop-agent
# Force the kernel to enable conntrack for both v4 and v6 by installing a
# tiny inert nftables table that references `ct`. The kernel only enables
# conntrack per-family on demand, so on a host with v4 NAT rules but no
# v6 `ct` references the IPv6 conntrack table stays empty (and qiftop
# can't show v6 flows). The `-` prefix makes this best-effort: if `nft`
# isn't installed we still start. ExecStopPost removes the table again.
ExecStartPre=-/usr/sbin/nft -f /usr/share/qiftop/qiftop-conntrack.nft
ExecStopPost=-/usr/sbin/nft delete table inet qiftop
User=root
# Capabilities the agent actually needs:
#   CAP_NET_ADMIN       — open libnetfilter_conntrack, enable nf_conntrack_acct,
#                         issue NETLINK_SOCK_DIAG dumps.
#   CAP_SYS_PTRACE      — read /proc/<pid>/fd/* symlinks for processes owned
#                         by other users. Without this, refreshProcWalk in
#                         SockDiagResolver silently sees EACCES on every
#                         non-root pid (root uid does NOT bypass
#                         ptrace_may_access()), the inode→pid map stays
#                         empty, and every flow reports pid=0 in the wire DTO
#                         even though `process-attribution-wire` is advertised.
#                         Also required to read /proc/<pid>/exe + /proc/<pid>/cwd
#                         symlinks served by the on-demand GetProcessDetails RPC.
#   CAP_DAC_READ_SEARCH — bypass DAC permission checks on /proc/<pid>/fd
#                         directories (mode 0500, owned by the process owner).
#                         A uid=0 process whose CapEff has been restricted via
#                         CapabilityBoundingSet does NOT automatically bypass
#                         DAC any more, so opendir("/proc/<other-uid-pid>/fd")
#                         returns EACCES without this cap. Symptom is the
#                         same as missing CAP_SYS_PTRACE: process attribution
#                         silently sees only its own-uid pids.
#   CAP_SYS_ADMIN       — needed by NetnsScanner's setns(2) into container
#                         network namespaces so flows held by container-side
#                         sockets attribute correctly. Without it, NetnsScanner
#                         fails to initialize and the netns-scan capability
#                         token is not advertised — process attribution still
#                         works for host-netns flows, just not for flows whose
#                         socket lives in a container netns.
CapabilityBoundingSet=CAP_NET_ADMIN CAP_SYS_PTRACE CAP_DAC_READ_SEARCH CAP_SYS_ADMIN
AmbientCapabilities=CAP_NET_ADMIN CAP_SYS_PTRACE CAP_DAC_READ_SEARCH CAP_SYS_ADMIN
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=no
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX AF_NETLINK AF_INET AF_INET6
# Allow ONLY network-namespace switching (deny user/mnt/pid/ipc/uts/cgroup).
# NetnsScanner setns(CLONE_NEWNET)s into each container's network namespace
# to dump its sockets and attribute container-side flows. RestrictNamespaces=yes
# (deny ALL namespace ops) silently blocked that setns — the scanner would
# enter the netns loop, fail setns with EPERM, and attribute every container
# flow as host. `=net` is the minimum that re-enables container attribution
# while still denying every other namespace type. (Also requires
# CAP_SYS_ADMIN, already granted above.)
RestrictNamespaces=net
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
Restart=on-failure
RestartSec=2

[Install]
# Activated on demand via DBus; no WantedBy needed by default. Admins who
# want it always on can `systemctl enable qiftop-agent.service`.
